# config/etc/supervisor/conf.d/crontab.conf
[program:cron]
command=/usr/sbin/cron -f
stopsignal=QUIT
autorestart=true