# config/etc/supervisor/conf.d/rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
stopsignal=QUIT
autorestart=true